# dot_config/autostart/blueman.desktop
[Desktop Entry]
Name=blueman
Exec=blueman-applet
StartupNotify=false
Terminal=false
Type=Application
