# dot_config/autostart/kanshi.desktop
[Desktop Entry]
Name=kanshi
Exec=kanshi
StartupNotify=false
Terminal=false
Type=Application
