# dot_config/autostart/swayosd.desktop
[Desktop Entry]
Name=swayosd
Exec=swayosd-server
StartupNotify=false
Terminal=false
Type=Application
