# dot_config/autostart/udiskie.desktop
[Desktop Entry]
Name=udiskie
Exec=udiskie --no-tray --no-appindicator --notify
StartupNotify=false
Terminal=false
Type=Application
